# crontab
# For host server
*/1 * * * * docker exec lightsquid_app_1 ./lightparser.pl today 2>> $PROJ_PATH/errors.log
0 4 * * * rm -f $PROJ_PATH/log/access.log
